# glint.service
[Unit]
Description=Control the glint service
After=mysql.service

[Service]
Type=simple
User=root
Group=root
EnvironmentFile=/usr/local/863/config/config.yaml
EnvironmentFile=/usr/local/863/config/xss.yaml
PIDFile=/run/glint.pid
ExecStart=/usr/local/863/bin/glint --websocket 127.0.0.1:3009 --cert /usr/local/863/certific/server.pem --key /usr/local/863/certific/server.key --dbconnect
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target